# etc/xdg/autostart/sxhkd.desktop
Name=sxhkd
Comment=HoyKey Daemon
Exec=/usr/bin/sxhkd
Terminal=false
Type=Application
